# wiki/WebcamPI.wiki
#summary One-sentence summary of this page.

= Introduction =

Einbindung der Raspberry PI CAM & Aktivierung eines Streamservers.

= Details =

1. Erstellung in /etc/init.d der Datei mjpg_streamer - http://wiki.ubuntuusers.de/MJPG-Streamer

2. Installation PI CAM -  http://www.raspberrypi.org/camera

3. Erstellung in /etc/init.d der Datei mjpg_stream

4. sudo chmod 0775 mjpg_stream

 # /mnt/vmtmpfs ersetzen durch das passende Verzeichnis (Ramdisc zu empfehlen!)
 # /home/pi/mjpeg/mjpg-streamer-code-182/mjpg-streamer/mjpg_streamer ersetzen durch das passende Verzeichnis je nachdem wo es installiert ist, Pfad notwendig beim Autostart!)

{{{
#!/bin/sh
# /etc/init.d/mjpg_stream

### BEGIN INIT INFO
# Provides:          mjpg_stream
# Required-Start:    $all
# Required-Stop:     $all
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: MJPG - Stream
# Description:       MJPG - Stream
### END INIT INFO

start()
{
  echo "Start mjpg-stream"
  /usr/bin/raspistill --nopreview -w 640 -h 480 -ISO 100 -ex auto -q 5 -o /mnt/vmtmpfs/pic.jpg -tl 1000 -t 9999999 -th 0:0:0 &
  /home/pi/mjpeg/mjpg-streamer-code-182/mjpg-streamer/mjpg_streamer -i "/usr/local/lib/input_file.so -f /mnt/vmtmpfs -n pic.jpg -d 5" -o "/usr/local/lib/output_http.so -w /usr/local/www" &
}

stop()
{
  echo "Stop mjpg-stream"
  kill -9 $(pidof mjpg_streamer) >/dev/null 2>&1
}

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart)
#    stop
    start
    ;;
  *)
    echo "Usage: $0 {start|stop|restart}"
    ;;
esac

exit 0
}}}

Danach in den Autostart laden in /etc/init.d:

{{{
sudo update-rd.d mjpg-streamer defaults
}}}

Einbindung als Bild In der GUI via Container:

HTML:   <img src="http://192.168.1.2:8080/?action=snapshot" width="320" height="240"/>